# access-control-web.service
[Unit]
Description=Access Control Web API Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Gira come root per avere tutti i permessi USB/hardware
User=root
Group=root
WorkingDirectory=/opt/access_control
Environment="PATH=/opt/access_control/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/access_control/src"

# Comando per avviare il servizio - SOLO web_api, tutto il resto è integrato
ExecStart=/opt/access_control/venv/bin/python /opt/access_control/src/api/web_api.py

# Con User=root non servono permessi extra, ha già tutto
PrivateDevices=no

# Restart automatico in caso di crash
StartLimitInterval=0
StartLimitBurst=10

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target