# access-control-web.service
[Unit]
Description=Access Control Web API Service
After=network.target

[Service]
Type=simple
User=wlkr42
Group=plugdev
WorkingDirectory=/opt/access_control
Environment="PATH=/opt/access_control/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/access_control/src"

# Comando per avviare il servizio
ExecStartPre=/opt/access_control/venv/bin/python /opt/access_control/src/hardware/crt285_wrapper.py
ExecStart=/opt/access_control/venv/bin/python /opt/access_control/src/api/web_api.py

# Permessi supplementari per accesso USB
SupplementaryGroups=dialout plugdev
PrivateDevices=no
DeviceAllow=/dev/bus/usb rw
DeviceAllow=/dev/ttyACM0 rw
DeviceAllow=/dev/ttyUSB0 rw

# Capacità per accesso diretto USB
AmbientCapabilities=CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_SYS_ADMIN

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target